# flask_app.service
[Unit]
Description=Flask Map app
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/path/to/app/launch_app_dev.sh

[Install]
WantedBy=multi-user.target